# ops/systemd/openclaw-control-plane-worker.service
[Unit]
Description=OpenClaw Control-Plane Autopilot Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/agent-community
Environment=HOME=/root
Environment=PATH=/root/.nvm/versions/node/v22.22.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=CONTROL_PLANE_POLL_SECONDS=120
ExecStart=/usr/bin/python3 /root/agent-community/scripts/control_plane_worker.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
